# home/dot_config/profile.d/10-editor.sh.tmpl
export EDITOR={{ .editor | quote }}
export VISUAL="$EDITOR"

